# idotmatrix-calendar-scheduler.service
[Unit]
Description=iDotMatrix Calendar Scheduler
Documentation=https://github.com/beny18241/python3-idotmatrix-client
After=network.target
Wants=network.target

[Service]
Type=simple
User=beny
Group=beny
WorkingDirectory=/home/beny/Desktop/idotmatrix
ExecStart=/home/beny/Desktop/idotmatrix/venv/bin/python /home/beny/Desktop/idotmatrix/calendar_scheduler_all.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=idotmatrix-calendar

# Environment variables
Environment=PATH=/home/beny/Desktop/idotmatrix/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/beny/Desktop/idotmatrix

# Security settings (relaxed for compatibility)
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target